# deploy/systemd/techpulse-refresh.timer
[Unit]
Description=Run TechPulse daily topic refresh at 08:05

[Timer]
OnCalendar=*-*-* 08:05:00
AccuracySec=1min
Persistent=true
Unit=techpulse-refresh.service

[Install]
WantedBy=timers.target
